# launchd/com.huapingyu.taskflow.local.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.huapingyu.taskflow.local</string>

    <key>WorkingDirectory</key>
    <string>/Users/huapingyu/dev/TaskFlow</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>npm run local:daemon</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/huapingyu/.nvm/versions/node/v24.14.1/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/huapingyu/dev/TaskFlow/.taskflow/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/huapingyu/dev/TaskFlow/.taskflow/launchd.log</string>
  </dict>
</plist>
